# webhook/deploy-webhook.service
[Unit]
Description=Deploy Webhook
After=network.target

[Service]
Type=simple
User=deploy
Group=deploy
WorkingDirectory=/opt/deploy-webhook
Environment=DEPLOY_WEBHOOK_CONFIG=/opt/deploy-webhook/config.json
Environment=DEPLOY_WEBHOOK_LOG=/var/log/deploy-webhook.log
Environment=PATH=/usr/local/bin:/usr/bin:/bin
ExecStart=/usr/bin/python3 /opt/deploy-webhook/deploy-webhook.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
